# templates/activemq.service.epp

[Unit]
Description=ActiveMQ
Wants=basic.target
After=basic.target network.target

[Service]
EnvironmentFile=-<%= $spath %>/activemq
Type=forking
PIDFile=/var/run/activemq.pid
ExecStart=<%= $home %>/bin/activemq start
ExecReload=<%= $home %>/bin/activemq reload
ExecStop=<%= $home %>/bin/activemq kill
User=<%= $user %>
Group=<%= $group %>

[Install]
WantedBy=multi-user.target
