# templates/activemq.service.epp

[Unit]
Description=ActiveMQ
Wants=basic.target
After=basic.target network.target

[Service]
EnvironmentFile=-<%= $spath %>/activemq
Type=forking
PIDFile="<%= $home %>/data/activemq.pid"
ExecStart=/opt/activemq/bin/activemq start
ExecReload=/opt/activemq/bin/activemq reload
ExecStop=/opt/activemq/bin/activemq kill
User=<%= $user %>
Group=<%= $group %>

[Install]
WantedBy=multi-user.target