# templates/activemq.service.epp

[Unit]
Description=ActiveMQ
Wants=basic.target
After=basic.target network.target

[Service]
EnvironmentFile=-<%= $spath %>/activemq
Type=forking
PIDFile="<%= $home %>/data/activemq.pid"
ExecStart="<%= $home %>/bin/start"
ExecStop="<%= $home %>/bin/stop"
ExecReload="<%= $home %>/bin/restart"
User=<%= $user %>
Group=<%= $group %>

[Install]
WantedBy=multi-user.target